# contrib/init/systemd/hbm.socket
[Unit]
Description=Harbormaster Plugin Socket

[Socket]
ListenStream=/run/docker/plugins/hbm.sock

[Install]
WantedBy=sockets.target
